# supervisor.conf
[supervisord]
nodaemon=true
[program:Mirrors-AutoSync]
command= python3 /Mirrors-AutoSync/Mirrors-AutoSync.py
